# server/status-email@.service
[Unit]
Description=Status email for %i

[Service]
Type=oneshot
ExecStart=/home/hmn/hmn/server/adminmailer.sh %i
User=root
Group=systemd-journal
